# scripts/systemd/x11vnc.service
[Unit]
Description=X11VNC Server for Virtual Display
After=xfce-desktop.service
Requires=xfce-desktop.service

[Service]
Type=simple
User=ubuntu
Environment="DISPLAY=:1"
Environment="XAUTHORITY=/home/ubuntu/.Xauthority"
ExecStart=/usr/bin/x11vnc -display :1 -forever -shared -rfbport 5900 -rfbauth /home/ubuntu/.vnc/passwd
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
